# vendors/SCS/common/inittab
klog:unknown:/sbin/klogd -n
http:unknown:/bin/boa -f /etc/boa.conf
intd:unknown:/bin/inetd
