# vendors/OpenGear/romfs/etc/default/inittab
console:unknown:/bin/agetty ttyAM0 115200
