# vendors/SecureEdge/romfs_x86/etc/default-static/inittab
boa::respawn:/bin/boa
dns::respawn:/bin/dnsmasq -d -l /etc/config/dhcpd.leases
dhsd::respawn:/bin/dhcpd eth0
