# vendors/SecureEdge/romfs_x86/etc/default/inittab
boa::respawn:/bin/boa
dns::respawn:/bin/dnsmasq -d
